# templates/default/shadowsocks.conf.erb
[program:<%= @name %>]
command=ssserver -c /etc/shadowsocks/<%= @name %>.json
autorestart=true
user=nobody
